# scripts/deploy/permit-api-deploy.service
[Unit]
Description=permit-api auto-deploy from origin/main
After=network-online.target

[Service]
Type=oneshot
User=will
# Invoked via bash: SELinux denies systemd direct-exec of scripts under /home.
ExecStart=/usr/bin/bash /home/will/permit-api-live/scripts/deploy/r730_deploy.sh
