# scripts/deploy/permit-api-deploy.timer
[Unit]
Description=permit-api auto-deploy poll (every 5 min)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
